# deployment/helpdesk-api.service
# /etc/systemd/system/helpdesk-api.service
# Salin ke VPS:
#   scp deployment/helpdesk-api.service root@<IP_VPS>:/etc/systemd/system/helpdesk-api.service
# Aktifkan:
#   systemctl daemon-reload
#   systemctl enable helpdesk-api
#   systemctl start helpdesk-api

[Unit]
Description=Helpdesk Chatbot FastAPI
After=network.target ollama.service

[Service]
Type=simple
User=chatbot
WorkingDirectory=/opt/helpdesk-backend
Environment="PATH=/opt/helpdesk-backend/venv/bin"
ExecStart=/opt/helpdesk-backend/venv/bin/uvicorn api_chatbot:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
